# confs/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq -d
autorestart=true
autostart=true
user=root
