# confs/supervisor/conf.d/tinc.conf
[program:tinc]
command=/usr/sbin/tincd -n twilight -D
autorestart=true
autostart=true
user=root