# stinc/supervisor/conf.d/chinadns.conf
[program:chinadns]
command=/usr/local/bin/chinadns -b 0.0.0.0 -c /etc/tinc/twilight/chnroute.txt -p 5353 -s 202.120.224.6,8.8.8.8
autorestart=true
directory=/
autostart=false
user=nobody
