# stinc/supervisor/conf.d/ss.conf
[program:ss]
command=ss-server -s 0.0.0.0 -p 8388 -m rc4-md5 -k testonly -t 60 -d 127.0.0.1 -u --fast-open
autorestart=true
directory=/usr/bin
autostart=false
user=root
